# server-postgresql/etc/crontab
##This should not be in the hot standby
 *  *  *  *  * root $SHELL /var/lib/pgsql/nodes_assurehealthy > 2>&1
