# server-pgpool/etc/crontab
  *  *  *  *  * root    /etc/pgpool-II/nodes_assurehealthy > 2>&1
